# deploy/systemd/reserve-auction-keeper-compose.service
[Unit]
Description=Ajna Reserve Auction Keeper (Docker Compose)
Requires=docker.service
After=docker.service network-online.target
Wants=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
WorkingDirectory=/opt/reserve-auction-keeper
ExecStart=/usr/bin/docker compose -f docker/docker-compose.yml up -d --build
ExecStop=/usr/bin/docker compose -f docker/docker-compose.yml stop
ExecReload=/usr/bin/docker compose -f docker/docker-compose.yml up -d --build
TimeoutStartSec=0
TimeoutStopSec=45

[Install]
WantedBy=multi-user.target
